# orchestrator/systemd/ujin-orchestrator.timer
[Unit]
Description=Fire the ujin orchestrator tick every 10 minutes
Requires=ujin-orchestrator.service

[Timer]
# First tick 5 minutes after boot/enable, then every 10 minutes (measured from the
# end of the previous oneshot run, so long builder ticks naturally space themselves).
OnBootSec=5min
OnUnitActiveSec=10min
# Catch up on missed ticks across machine-off windows.
Persistent=true
Unit=ujin-orchestrator.service

[Install]
WantedBy=timers.target
